# recipes-quectel/quectel_cm/files/lte.service
[Unit]
Description=LTE service
Before=network.service

[Service]
ExecStart=/usr/sbin/quectel_CM -s ctnet

[Install]
WantedBy=default.target
